# frontend/supervisord.conf
[supervisord]
nodaemon=true

#[program:pm2 ]
#command=pm2 start /var/www/html/SHOP/process.yml

[program:apache_restart]
command=service apache2 restart 
